# systemd/chapel-bells-web.service
[Unit]
Description=Midway UMC Bells - Web Dashboard
After=network.target pipewire.service pipewire-pulse.service
Wants=pipewire.service pipewire-pulse.service

[Service]
Type=simple
User=bells
Group=bells
WorkingDirectory=/opt/bells

ExecStart=/opt/bells/venv/bin/python run_web_ui.py

Restart=on-failure
RestartSec=10
StartLimitIntervalSec=300
StartLimitBurst=5

# ── Security hardening ──
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectKernelLogs=yes
ProtectControlGroups=yes
ProtectClock=yes
ProtectHostname=yes
RestrictSUIDSGID=yes
RestrictRealtime=yes
RestrictNamespaces=yes
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
LockPersonality=yes
SystemCallFilter=@system-service
SystemCallArchitectures=native

# Allow access to audio devices
PrivateDevices=no
SupplementaryGroups=audio

ReadWritePaths=/opt/bells /var/log/bells

# Environment
Environment="PYTHONUNBUFFERED=1"
Environment="SDL_AUDIODRIVER=alsa"

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=midway-bells-web

# Resource limits
MemoryMax=128M
CPUQuota=25%

[Install]
WantedBy=multi-user.target
