# systemd/chapel-bells.service
[Unit]
Description=Midway UMC Bells - Automated Church Bell System
Documentation=https://github.com/HeyJonesR/BellSystem
After=network.target time-sync.target sound.target
Wants=network-online.target

[Service]
Type=simple
User=bells
Group=bells
WorkingDirectory=/opt/bells

ExecStart=/opt/bells/venv/bin/python -m chapel_bells --config /opt/bells/config/schedule.json --log-file /var/log/bells/bells.log

# Auto-restart on failure
Restart=on-failure
RestartSec=10
StartLimitIntervalSec=300
StartLimitBurst=5

# ── Security hardening ──
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectKernelLogs=yes
ProtectControlGroups=yes
ProtectClock=yes
ProtectHostname=yes
RestrictSUIDSGID=yes
RestrictRealtime=yes
RestrictNamespaces=yes
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
LockPersonality=yes
SystemCallFilter=@system-service
SystemCallArchitectures=native

# Allow access to audio devices (sound board / DAC HAT)
PrivateDevices=no
SupplementaryGroups=audio

# Writable paths
ReadWritePaths=/opt/bells /var/log/bells

# Resource limits
MemoryMax=256M
CPUQuota=50%

# Environment
Environment="PYTHONUNBUFFERED=1"
StandardOutput=journal
StandardError=journal
SyslogIdentifier=midway-bells

KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=10

[Install]
WantedBy=multi-user.target
